# scripts/ai.openclaw.tailscale-serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.tailscale-serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Tailscale.app/Contents/MacOS/Tailscale</string>
        <string>serve</string>
        <string>--bg</string>
        <string>http://127.0.0.1:18789</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/openclaw/tailscale-serve.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/tailscale-serve.log</string>
</dict>
</plist>
